# deploy/scammer-bot.service
[Unit]
Description=Scammer List Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/scammer-list-bot
EnvironmentFile=/home/botuser/scammer-list-bot/.env
ExecStart=/home/botuser/scammer-list-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
